# devops/comandos_e_scripts_linux/crontab
* * * * * /home/caiohbs/monitoramento.sh >> /home/caiohbs/saida_nginx.txt